# qqbot.service
[Unit]
Description=QQ Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/TableHelper
ExecStart=/usr/bin/python3 /root/TableHelper/run_qq_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/qqbot.log
StandardError=append:/var/log/qqbot.log

[Install]
WantedBy=multi-user.target
